# deploy/liberty-basketball-analysis.service
[Unit]
Description=Liberty Basketball Analysis
After=network.target

[Service]
Type=simple
User=REPLACE_WITH_SERVER_USER
WorkingDirectory=/opt/liberty-basketball-analysis
Environment=LIBERTY_DEBUG=0
Environment=PORT=8080
Environment=LIBERTY_DATABASE=/opt/liberty-basketball-analysis/film_analysis.db
Environment=LIBERTY_UPLOAD_FOLDER=/opt/liberty-basketball-analysis/uploads
ExecStart=/opt/liberty-basketball-analysis/.venv/bin/python /opt/liberty-basketball-analysis/app.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
